# sparkgpt.service
# TEMPLATE - Not installed. Copy to ~/.config/systemd/user/ and run:
#   systemctl --user daemon-reload
#   systemctl --user enable sparkgpt
#   systemctl --user start sparkgpt

[Unit]
Description=ClawChat Voice Server
After=network.target

[Service]
Type=simple
User=heisenberg
WorkingDirectory=/home/heisenberg/clawd/sparkgpt
ExecStart=/usr/bin/node src/server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
